# package/infineon-utilities/feeds/ifx_base_files/target-ltqcpe-base-files/files.all/etc/init.d/source_modelconfig.sh
#!/bin/sh /etc/rc.common
#START=03
start() {
	if [ -r /etc/rc.d/model_config.sh ]; then
		. /etc/rc.d/model_config.sh
	fi
}
