# package/infineon-utilities/feeds/ifx_base_files/ifx-base-files-common/files.all/etc/init.d/ifx_show_version.sh
#!/bin/sh /etc/rc.common
# Copyright (C) 2010 OpenWrt.org
# Copyright (C) 2010 Lantiq.com

START=99

bindir=/opt/lantiq/bin

start() {
   # start the dsl daemon
   [ -e ${bindir}/show_version.sh ] && ${bindir}/show_version.sh
	[ -e /opt/lantiq/image_version ] && cat /opt/lantiq/image_version
}
