# package/infineon-utilities/feeds/ifx_base_files/base-files-vr9-bsp/files.all/etc/init.d/demo_setup
#!/bin/sh /etc/rc.common

START=98

start() {
cd /flash/BSP-Test-VR9/
./bringup_xdsl_mode.sh
}
