# package/infineon-utilities/feeds/ifx_base_files/base-files-vr9-bsp/files.all/etc/init.d/nfext.sh
#!/bin/sh /etc/rc.common

#START=25

start() {
	kernel_version=`uname -r`
	/sbin/insmod /lib/modules/${kernel_version}/ifx_nfext_core.ko
	/sbin/insmod /lib/modules/${kernel_version}/ifx_nfext_ppp.ko
}
