# package/infineon-utilities/feeds/ifx_base_files/base-files-vr9-bsp/files.all/etc/init.d/wlan_bringup.sh
#!/bin/sh /etc/rc.common
#START=45
start() {
	if [ "$CONFIG_FEATURE_IFX_WIRELESS" = "1" ]; then
		/etc/rc.d/rc.bringup_wlan load
	fi
}
